# templates/vector.service.j2
[Unit]
Description=Vector Log Agent
After=network.target

[Service]
Type=simple
User=vector
Group=vector
ExecStart={{ vector_install_dir }}/bin/vector --config {{ vector_config_file }}
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ vector_config_dir }}

[Install]
WantedBy=multi-user.target